# machines/work-desktop/systemd/user/xcmenu.service
[Unit]
Description=X clipboard menu
After=xorg.service

[Service]
ExecStart=/usr/bin/xcmenu

[Install]
WantedBy=graphical.target
